# blindwatermark/docker/build/conf/service.conf
[program:srp]
command=/var/opt/bin/startreporter.sh
environment=IP=%(ENV_IP)s,PORT=%(ENV_PORT)s,SERVICE_NAME=%(ENV_SERVICE_NAME)s,CONSUL_ADDR=%(ENV_CONSUL_ADDR)s
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
[program:httpGateway]
command=/app/golangexecapi -addr=:%(ENV_PORT)s -accessId="*****" -accessSecret="****" -bucket=yuanli-public -endPoint=oss-cn-beijing-internal.aliyuncs.com
environment=IP=%(ENV_IP)s,PORT=%(ENV_PORT)s,SERVICE_NAME=%(ENV_SERVICE_NAME)s,CONSUL_ADDR=%(ENV_CONSUL_ADDR)s
autostart=true
autorestart=true
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/data/output.log
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
[program:gtail]
command=/var/opt/bin/gtail-linux-64-1.0.0  -key "ky-17278a3ff9****" -secret "sec-e73eb65******" -sid "0e6f2****" /data/output.log
environment=IP=%(ENV_IP)s,PORT=%(ENV_PORT)s,SERVICE_NAME=%(ENV_SERVICE_NAME)s,CONSUL_ADDR=%(ENV_CONSUL_ADDR)s
autostart=true
autorestart=true
directory = /data/
priority=10
stdout_events_enabled=true
stderr_events_enabled=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
